# docker-postfix/supervisord.conf
[supervisord]
nodaemon=true

[program:syslog]
command=/usr/sbin/rsyslogd -n

[program:postfix]
command=/usr/local/bin/start_postfix.sh

